# deploy/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn Diginnocent.wsgi:application --bind 0.0.0.0:8000 --workers 4
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stopsignal=TERM

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
